# mcp-deployer.service
[Unit]
Description=Linux Application Deployer MCP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mcp-deployer
WorkingDirectory=/opt/mcp/deployer
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8001 --log-config logging_config.py
Restart=on-failure
RestartSec=5s

# Environment variables for logging
EnvironmentFile=-/etc/default/mcp-deployer
Environment="MCP_LOG_LEVEL=INFO"
Environment="MCP_LOGS_DIR=/var/log/mcp-linux-deployer"
Environment="MCP_LOG_MAX_BYTES=10485760"
Environment="MCP_LOG_BACKUP_COUNT=10"
Environment="MCP_LOG_FORMAT=detailed"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/mcp /var/log/mcp-linux-deployer

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target

