# prod-deployment-scripts/mcp-linux-app-deployer.service
[Unit]
Description=Linux Application Deployer MCP Server
After=network.target nginx.service
Requires=nginx.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mcp/linux-app-deployer

Environment=PYTHONUNBUFFERED=1
Environment=UV_CACHE_DIR=/opt/mcp/linux-app-deployer/.uv-cache

ExecStart=/opt/mcp/linux-app-deployer/.venv/bin/python server.py

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
