# prod-deployment-scripts/mcp-linux-app-deployer.service
[Unit]
Description=Linux Application Deployer MCP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mcp/linux-app-deployer
ExecStart=uv run python /opt/mcp/linux-app-deployer/server.py
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/mcp /var/log/mcp-linux-deployer

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target

